# weatherstation.service
[Unit]
Description = Weather Display Station
After = network.target
 
[Service]
Type = simple
ExecStart = <APP_DIR>/weather-station.sh start /opt/weather_station_cfg/config.cfg
User = root
 
[Install]
WantedBy = multi-user.target